# deploy/receipt-ranger.service
[Unit]
Description=Receipt Ranger Streamlit Application
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/receipt-ranger
Environment="PATH=/home/ubuntu/receipt-ranger/venv/bin"
Environment="ENABLE_GOOGLE_SHEETS=true"
EnvironmentFile=/home/ubuntu/receipt-ranger/.env
ExecStart=/home/ubuntu/receipt-ranger/venv/bin/streamlit run app.py --server.headless true
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/ubuntu/receipt-ranger/logs/streamlit.log
StandardError=append:/home/ubuntu/receipt-ranger/logs/streamlit.error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/receipt-ranger/output /home/ubuntu/receipt-ranger/logs /home/ubuntu/receipt-ranger/processed_receipts.json

[Install]
WantedBy=multi-user.target
